# units/netdumper.service
[Unit]
Description=netdumper
After=sshdgenkeys.service

[Service]
Type=notify
ExecStart=/home/nekit/bin/netdumper