# units/netdumper.service
[Unit]
Description=netdumper

[Service]
Type=notify
ExecStart=/home/nekit/bin/netdumper
